# deb/toplogger/usr/lib/systemd/system/toplogger.service
[Unit]
Description=top logger per-minute

[Service]
ExecStart=/usr/lib/toplogger/toplogger.sh

[Install]
WantedBy=network.target