# deb/toplogger/DEBIAN/prerm
#!/bin/bash
set -e

# Service
systemctl stop toplogger
systemctl disable toplogger

# AppArmor
aa-disable /etc/apparmor.d/usr.lib.toplogger.toplogger.sh
