# deb/toplogger/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop toplogger
systemctl disable toplogger
